# tiny.service
[Unit]
Description=DataPAI × TinyFish Stock Radar — tinyfish.datap.ai
After=network.target

[Service]
# Long-running Next.js server (use 'simple', NOT 'oneshot')
Type=simple
User=ec2-user
Group=ec2-user

# Load API keys (TINYFISH_API_KEY etc.)
EnvironmentFile=/home/ec2-user/.env.dev

# Port — Next.js reads PORT automatically; no code change needed
Environment=PORT=3085
Environment=NODE_ENV=production
Environment=HOSTNAME=0.0.0.0
# Explicitly add Node 20 bin to PATH so npm shell scripts can find node
Environment=PATH=/home/ec2-user/.nvm/versions/node/v20.20.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

WorkingDirectory=/home/ec2-user/git/datapai-tinyfish

# Call node directly — more reliable than npm in systemd (no shell PATH dependency)
ExecStart=/home/ec2-user/.nvm/versions/node/v20.20.0/bin/node node_modules/next/dist/bin/next start

# Restart on crash, wait 10s between attempts
Restart=always
RestartSec=10

# Logging — view with: journalctl -u tiny -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=datapai-tinyfish

[Install]
WantedBy=multi-user.target
